# systemd/idrac_fancontrol.timer
[Unit]
Description=Run IDRAC Fan Control every 5 minutes

[Timer]
OnCalendar=*:0/5
Persistent=false

[Install]
WantedBy=timers.target
